# deploy/label-print-server.service
[Unit]
Description=Label Print Server (WPS One-Click Print)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/www/label-print-server
ExecStart=/usr/bin/python3 -m uvicorn web_app:app --host 0.0.0.0 --port 8090
Restart=always
RestartSec=5
StandardOutput=append:/www/label-print-server/data/logs/server.log
StandardError=append:/www/label-print-server/data/logs/server_error.log

[Install]
WantedBy=multi-user.target
